# publisher.service
# Machine Status Publisher Service
[Unit]
Description=Machine Status MQTT Publisher
After=network.target mosquitto.service

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /usr/local/bin/machine-status-publisher.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target