# subscriber.service
[Unit]
Description=Machine Status MQTT Subscriber
After=network.target mosquitto.service postgresql.service

[Service]
Type=simple
User=root
EnvironmentFile=/etc/machine-status/mqtt.env
EnvironmentFile=/etc/machine-status/db.env
ExecStart=/opt/machine-status/venv/bin/python3 /opt/machine-status/subscriber/machine_status_subscriber.py
Restart=on-failure
RestartSec=10
WorkingDirectory=/opt/machine-status
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
EOL